# deploy/linux/cc-connect-qq.service
[Unit]
Description=cc-connect QQ bot
After=network-online.target onebot-group-proxy.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/chatbot-qq
EnvironmentFile=-/etc/chatbot-qq.env
Environment=HOME=/root/.codex-qq-home
Environment=USER=root
ExecStartPre=/bin/bash /opt/chatbot-qq/deploy/linux/wait-onebot-ports.sh
ExecStart=/usr/bin/env cc-connect --config /root/.cc-connect-qq/config.toml --force
Restart=always
RestartSec=10
# Minimum required for Codex/bubblewrap command execution.
# Do not change to NoNewPrivileges=true, add CapabilityBoundingSet=, or remove
# AF_NETLINK unless `codex sandbox linux -- /bin/sh -lc 'pwd; ls -ld .'`
# still passes under this unit's systemd sandbox.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/chatbot-qq /root/.cc-connect-qq /root/.codex-qq-home /var/log
ReadOnlyPaths=/opt/chatbot-qq/AGENTS.md /opt/chatbot-qq/docs /opt/chatbot-qq/configs
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
SystemCallArchitectures=native
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077
LimitNOFILE=8192
StandardOutput=append:/var/log/cc-connect-qq.log
StandardError=append:/var/log/cc-connect-qq.log

[Install]
WantedBy=multi-user.target
