# deploy/linux/chatbot-qq-profile-update.service
[Unit]
Description=Update CHATBOT QQ user profiles from recent chat records
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/chatbot-qq
EnvironmentFile=-/etc/chatbot-qq.env
Environment=HOME=/root/.codex-qq-home
Environment=USER=root
ExecStart=/bin/bash /opt/chatbot-qq/scripts/update-user-profiles.sh --all
TimeoutStartSec=45min
# Minimum required for Codex/bubblewrap command execution.
# Do not change to NoNewPrivileges=true, add CapabilityBoundingSet=, or remove
# AF_NETLINK unless `codex sandbox linux -- /bin/sh -lc 'pwd; ls -ld .'`
# still passes under this unit's systemd sandbox.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/opt/chatbot-qq /root/.codex-qq-home /var/log
ReadOnlyPaths=/opt/chatbot-qq/AGENTS.md /opt/chatbot-qq/docs /opt/chatbot-qq/configs /opt/chatbot-qq/scripts
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
SystemCallArchitectures=native
RestrictSUIDSGID=true
LockPersonality=true
UMask=0077
StandardOutput=append:/var/log/chatbot-qq-profile-update.log
StandardError=append:/var/log/chatbot-qq-profile-update.log
